# srcs/influxdb/supervisor-telegraf.conf
[supervisord]
nodaemon=true

[program:influxdb]
command=sh -c "influxd -config /etc/influxdb.conf && kill -s SIGTERM $(cat supervisord.pid)"

[program:telegraf]
command=sh -c "/usr/bin/telegraf -config /etc/telegraf/telegraf.conf && kill -s SIGTERM $(cat supervisord.pid)"